# vllm/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:vllm]
command=python3 -m vllm.entrypoints.openai.api_server --host 0.0.0.0 --port 8000 --model %(ENV_VLLM_MODEL)s --gpu-memory-utilization %(ENV_VLLM_GPU_MEM_UTIL)s --max-model-len %(ENV_VLLM_MAX_MODEL_LEN)s --enable-prefix-caching
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true